# deploy/systemd/pitgun-telemetryd.service
[Unit]
Description=Pitgun telemetry ingestion daemon (behind Caddy)
After=network.target
ConditionPathExists=/opt/pitgun/pitgun/target/release/pitgun-telemetryd

[Service]
Type=simple
WorkingDirectory=/opt/pitgun/pitgun
ExecStart=/opt/pitgun/pitgun/target/release/pitgun-telemetryd
Restart=on-failure
RestartSec=2s
# Adjust the runtime user to match your deployment.
User=pitgun
Group=pitgun
Environment=PITGUN_TELEMETRY_BIND=127.0.0.1:8080
Environment=PITGUN_TELEMETRY_DATA_DIR=/opt/pitgun/telemetry/data
# Example logging level (uses tracing-subscriber env filter)
Environment=RUST_LOG=info

[Install]
WantedBy=multi-user.target

# Notes:
# - TLS is terminated by Caddy; keep PITGUN_TELEMETRY_BIND on localhost.
# - Ensure /opt/pitgun/telemetry/data is writable by the service user.
# - Repo is expected at /opt/pitgun/pitgun with release binary in target/release/.
